# respaldo-fotos-familia.service
[Unit]
Description=Respaldo de fotos familiares DiazSantaM a Google Drive
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pedro
ExecStart=/home/pedro/respaldo-fotos-familia.sh
Restart=always
RestartSec=30
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
